# .config/systemd/user/peapix-wallpaper-archive-download.service
[Unit]
Description=Fetch archive

[Service]
Type=oneshot
WorkingDirectory=/volumes/media/photo/.stage/peapix
ExecStart=/usr/bin/bash -c 'source "%h/.local/lib/relocatable-perl-dev.rc"; "%h/.local/scripts-ext/peapix/peapix-wallpaper-archive-download"'

[Install]
WantedBy=default.target
